# launchd/templates/com.myndaix.claw-cursor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myndaix.claw-cursor</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Applications/ClawCursor.app/Contents/MacOS/claw-cursor-engine</string>
        <string>--socket</string>
    </array>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claw-cursor/logs/engine.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAW_CONFIG_PATH</key>
        <string>__HOME__/.claw-cursor</string>
    </dict>
</dict>
</plist>
